# docker/supervisord/workers/test-worker.ini
[program:test-worker]
process_name=%(program_name)s_%(process_num)02d
command=php /backend/artisan queue:work
autostart=true
autorestart=true
stopasgroup=true
killasgroup=true
numprocs=1
redirect_stderr=true
stderr_logfile=/var/log/test-worker.log
stdout_logfile=/var/log/test-worker.log
stopwaitsecs=-1
